# roles/tgbot/templates/tgfeedbackbotworker.service
[Unit]
Description=Telegram Feedback Bot Service
Documentation="Ooops, no docs"
After=network.target

[Service]
Type=simple
User={{ os_user_name }}
ExecStart=/usr/bin/python3 /home/{{ os_user_name }}/tgfeedbackbotworker/main.py
StandardOutput=append:/var/log/tgfeedbackbotworker/out.log
StandardError=append:/var/log/tgfeedbackbotworker/error.log
Restart=on-failure

[Install]
WantedBy=multi-user.target
